# raspberry-pi/systemd/can-ids.service
[Unit]
Description=CAN-IDS - Controller Area Network Intrusion Detection System
Documentation=https://github.com/yourusername/can-ids
After=network.target multi-user.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/can-ids
Environment="PATH=/home/pi/can-ids/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/bin/sleep 10
ExecStartPre=/sbin/ip link set can0 type can bitrate 500000
ExecStartPre=/sbin/ip link set up can0
ExecStart=/home/pi/can-ids/venv/bin/python /home/pi/can-ids/main.py -i can0 --config /home/pi/can-ids/config/can_ids_rpi4.conf
ExecStop=/sbin/ip link set down can0
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=can-ids

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/can-ids/logs /home/pi/can-ids/data

# Resource limits
MemoryLimit=512M
CPUQuota=70%

[Install]
WantedBy=multi-user.target